# optional/dropbox/files/systemd/rclone-dropbox.service
[Unit]
Description=Mirror Dropbox to %h/Dropbox via rclone sync (read-only, one-way)
Documentation=https://rclone.org/commands/rclone_sync/
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone sync dropbox: %h/Dropbox --max-delete 50
Nice=10
